# noip-automatic-update/rc.local
#!/bin/sh -e
'/usr/local/bin/noip2'
exit 0
